# pam-limiter-cleanup.service
[Unit]
Description=PAM Limiter Cleanup Service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/pam-limiter-cleanup.sh
